# crontab
0 5 * * * /path/to/repo/run.sh >> /path/to/log
0 8 * * * /path/to/repo/render_widget.sh >> /path/to/log
